# production/systemd/claude-code-enforcement.service
[Unit]
Description=Claude Code Enforcement System
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=devflow
Group=devflow
WorkingDirectory=/opt/devflow/enforcement
ExecStart=/usr/bin/node /opt/devflow/enforcement/dist/bootstrap/claude-code-bootstrap.js
Restart=always
RestartSec=10
Environment=NODE_ENV=production
Environment=LOG_LEVEL=info
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target